# securedrop/debian/config/lib/systemd/system/securedrop-cleanup-ossec.service
[Unit]
Description=Cleanup OSSEC diff queue

[Service]
Type=oneshot
ExecStart=/usr/bin/securedrop-cleanup-ossec.py
User=root
