# securedrop/debian/config/lib/systemd/system/securedrop-reboot-required.service
[Unit]
Description=Touch /var/run/reboot-required file

[Service]
Type=oneshot
ExecStart=/usr/bin/touch /var/run/reboot-required
User=root
