# securedrop/debian/config/lib/systemd/system/securedrop-remove-packages.timer
[Unit]
Description=Remove ufw and haveged if installed

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
